# photoloop.service
# Copyright (c) 2025 Luc Vincent. All Rights Reserved.
[Unit]
Description=PhotoLoop - Raspberry Pi Photo Frame
Documentation=https://github.com/lucvincent/photoloop
After=network.target graphical.target
Wants=graphical.target

[Service]
Type=notify
User=luc
Group=luc

# Working directory
WorkingDirectory=/opt/photoloop

# Environment - Wayland/labwc on Raspberry Pi OS Bookworm
# Note: Assumes user has UID 1000 (default for first user on Pi)
Environment=XDG_RUNTIME_DIR=/run/user/1000
Environment=WAYLAND_DISPLAY=wayland-0
Environment=SDL_VIDEODRIVER=wayland
Environment=PHOTOLOOP_LOG_DIR=/var/log/photoloop
Environment=PHOTOLOOP_API_URL=http://localhost:8080

# Main executable
ExecStart=/opt/photoloop/venv/bin/python -m photoloop.src.main -c /etc/photoloop/config.yaml

# Restart policy
Restart=on-failure
RestartSec=5

# Watchdog: systemd restarts service if no WATCHDOG=1 within 60s
WatchdogSec=60

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=photoloop

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=true

# Allow write access to required directories
ReadWritePaths=/var/lib/photoloop
ReadWritePaths=/var/log/photoloop
ReadWritePaths=/etc/photoloop

[Install]
WantedBy=graphical.target
